# supervisord.conf
[supervisord]

[program:otbox.py]

command=python /home/inria/Documents/opentestbed/otbox.py
autostart=true
autorestart=ture
